# scripts/distribution/com.finalchapter.digest.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.finalchapter.digest</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/qili/git/finalchapter/scripts/distribution/daily.sh</string>
  </array>

  <!-- Run every day at 08:00 local. If the Mac is asleep at 8, launchd runs
       it at the next wake. -->
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key><integer>8</integer>
    <key>Minute</key><integer>0</integer>
  </dict>

  <key>RunAtLoad</key>
  <false/>

  <key>StandardOutPath</key>
  <string>/Users/qili/finalchapter-digests/launchd.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/qili/finalchapter-digests/launchd.err.log</string>
</dict>
</plist>
